# guest/network-manager.conf
# network-manager - network connection manager
#
# The Network Manager daemon manages the system's network connections,
# automatically switching between the best available.

description	"network connection manager"

start on (local-filesystems
	  and started dbus)
stop on stopping dbus

pre-start script
   VM_NUM=`/sbin/hc_maker`
   if [ ${VM_NUM} != -1 ]; then
      # The following should perhaps be put in a separate script
      ifconfig eth0 192.168.0.$((100+${VM_NUM})) netmask 255.255.0.0
      route add default gw 192.168.0.2
      # Set up DNS
      echo "nameserver YOUR_NAMESERVER
domain YOUR_DOMAIN
search YOUR_DOMAIN" >/etc/resolv.conf
      stop
   fi
end script

expect fork
respawn

exec NetworkManager
